# overclock.service
[Unit]
Description=Increase GPU core and memory clocks

[Service]
Type=oneshot
ExecStart=/usr/bin/overclock

[Install]
WantedBy=multi-user.target
